# appliance/packer/files/sixtyops-openrc
#!/sbin/openrc-run
# OpenRC service for SixtyOps Manager

name="sixtyops"
description="SixtyOps Manager"
DATA_MODE_FILE="/etc/sixtyops/data-mode"

depend() {
    need docker
    after docker
}

start() {
    ebegin "Starting ${name}"
    # Wait up to 60 seconds for Docker socket to be available
    local i=0
    while [ $i -lt 60 ]; do
        if docker info >/dev/null 2>&1; then
            break
        fi
        [ $((i % 10)) -eq 0 ] && [ $i -gt 0 ] && einfo "Waiting for Docker... (${i}s)"
        sleep 1
        i=$((i + 1))
    done
    if [ $i -ge 60 ]; then
        eerror "Docker not ready after 60 seconds"
        eend 1
        return 1
    fi
    # Verify /data state before starting containers (prevents silent data loss)
    local data_mode
    data_mode=$(cat "$DATA_MODE_FILE" 2>/dev/null)
    if [ -z "$data_mode" ]; then
        if grep -q 'sixtyops-data' /etc/fstab 2>/dev/null; then
            data_mode="partition"
        else
            data_mode="rootfs"
        fi
    fi
    if [ "$data_mode" = "partition" ]; then
        if ! mountpoint -q /data; then
            eerror "/data partition is not mounted — refusing to start (data would be lost)"
            eend 1
            return 1
        fi
    elif ! mountpoint -q /data; then
        einfo "/data is configured on root filesystem mode"
    fi

    if ! touch /data/.write-test 2>/dev/null; then
        eerror "/data is not writable — refusing to start"
        eend 1
        return 1
    fi
    rm -f /data/.write-test

    # Ensure self-signed SSL cert exists for nginx (nginx will crash-loop without it)
    if [ ! -f /data/certs/selfsigned.crt ] || [ ! -f /data/certs/selfsigned.key ]; then
        einfo "Generating self-signed SSL certificate..."
        mkdir -p /data/certs
        if ! openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
            -keyout /data/certs/selfsigned.key \
            -out /data/certs/selfsigned.crt \
            -subj "/C=US/ST=State/L=City/O=SixtyOps/CN=localhost" \
            2>/dev/null; then
            eerror "Failed to generate SSL certificate — nginx cannot start without it"
            eend 1
            return 1
        fi
    fi

    cd /opt/sixtyops
    docker compose up -d
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    cd /opt/sixtyops
    timeout 90 docker compose down 2>/dev/null
    eend $?
}

status() {
    if timeout 5 docker ps --filter name=sixtyops-management --filter status=running -q 2>/dev/null | grep -q .; then
        einfo "service started"
        return 0
    else
        einfo "service stopped"
        return 3
    fi
}
